# extras/borgmatic-borg-backup/etc/cron.d/borgmatic
# run borgmatic hourly

0 * * * * root borgmatic --syslog-verbosity 1
